# postiz/rootfs/etc/s6-overlay/s6-rc.d/redis/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Redis finish script
# ==============================================================================
bashio::log.info "Redis stopped."
