# postiz/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv bashio
exec redis-server --port 6380 --dir /data/redis --appendonly yes --protected-mode no --bind 127.0.0.1
