# postiz/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Redis service
# ==============================================================================
bashio::log.info "Starting Redis..."

exec su-exec redis redis-server /etc/redis.conf
