# postiz/rootfs/etc/s6-overlay/s6-rc.d/postiz/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Postiz finish script
# ==============================================================================
bashio::log.info "Postiz stopped."
pm2 delete all 2>/dev/null || true
