# postiz/rootfs/etc/s6-overlay/s6-rc.d/postiz/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Postiz application service (nginx + pm2)
# ==============================================================================

# Load environment
set -a
# shellcheck source=/dev/null
source /etc/postiz-env
set +a

# Wait for PostgreSQL to be ready
bashio::log.info "Waiting for PostgreSQL..."
for i in $(seq 1 30); do
    if /usr/lib/postgresql/17/bin/pg_isready -h 127.0.0.1 -p 5432 -U postiz -q 2>/dev/null; then
        break
    fi
    if [ "$i" -eq 30 ]; then
        bashio::log.error "PostgreSQL did not become ready in time!"
        exit 1
    fi
    sleep 1
done
bashio::log.info "PostgreSQL is ready."

# Wait for Redis to be ready
bashio::log.info "Waiting for Redis..."
for i in $(seq 1 30); do
    if redis-cli -h 127.0.0.1 ping 2>/dev/null | grep -q PONG; then
        break
    fi
    if [ "$i" -eq 30 ]; then
        bashio::log.error "Redis did not become ready in time!"
        exit 1
    fi
    sleep 1
done
bashio::log.info "Redis is ready."

cd /app || exit 1

# Run Prisma DB migrations
bashio::log.info "Running database migrations..."
pnpm run prisma-db-push 2>&1 || bashio::log.warning "prisma-db-push failed or not needed"

# Start nginx + pm2 (matches the official image CMD)
bashio::log.info "Starting nginx and pm2..."
nginx && exec pnpm run pm2
