# postiz/rootfs/etc/s6-overlay/s6-rc.d/postiz/run
#!/command/with-contenv bashio

set -a
source /etc/postiz-env
set +a

# Wait for PostgreSQL
bashio::log.info "Waiting for PostgreSQL..."
for i in $(seq 1 60); do
    if /usr/lib/postgresql/17/bin/pg_isready -h /run/postgresql -p 5433 > /dev/null 2>&1; then
        bashio::log.info "PostgreSQL is ready."
        break
    fi
    [ "$i" -eq 60 ] && { bashio::log.error "PostgreSQL timeout"; exit 1; }
    sleep 1
done

# Create user/db if first run
if [ -f /data/.needs_db_setup ]; then
    bashio::log.info "Creating database and user..."
    /usr/lib/postgresql/17/bin/psql -h /run/postgresql -p 5433 -U postgres -c "CREATE USER postiz WITH PASSWORD 'postiz-password' SUPERUSER;" 2>&1 || true
    /usr/lib/postgresql/17/bin/psql -h /run/postgresql -p 5433 -U postgres -c "CREATE DATABASE \"postiz-db\" OWNER postiz;" 2>&1 || true
    rm -f /data/.needs_db_setup
    bashio::log.info "Database created."
fi

# Wait for Redis
bashio::log.info "Waiting for Redis..."
for i in $(seq 1 30); do
    if redis-cli -p 6380 ping > /dev/null 2>&1; then
        bashio::log.info "Redis is ready."
        break
    fi
    sleep 1
done

cd /app

# Run DB migrations
bashio::log.info "Running Prisma migrations..."
pnpm run prisma-db-push 2>&1

# Start nginx + pm2
bashio::log.info "Starting Postiz..."
nginx
exec pnpm run pm2
