# ha-addon/rootfs/etc/s6-overlay/s6-rc.d/postiz/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Postiz finish script
# ==============================================================================
bashio::log.info "Postiz stopped."
