# ha-addon/rootfs/etc/s6-overlay/s6-rc.d/postiz/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Postiz application service
# ==============================================================================
bashio::log.info "Starting Postiz..."

# Load environment
set -a
# shellcheck source=/dev/null
source /etc/postiz.env
set +a

# Wait for PostgreSQL to be ready
bashio::log.info "Waiting for PostgreSQL..."
for i in $(seq 1 30); do
    if pg_isready -h 127.0.0.1 -p 5432 -U postiz -q 2>/dev/null; then
        break
    fi
    if [ "$i" -eq 30 ]; then
        bashio::log.error "PostgreSQL did not become ready in time!"
        exit 1
    fi
    sleep 1
done
bashio::log.info "PostgreSQL is ready."

# Wait for Redis to be ready
bashio::log.info "Waiting for Redis..."
for i in $(seq 1 30); do
    if redis-cli -h 127.0.0.1 ping 2>/dev/null | grep -q PONG; then
        break
    fi
    if [ "$i" -eq 30 ]; then
        bashio::log.error "Redis did not become ready in time!"
        exit 1
    fi
    sleep 1
done
bashio::log.info "Redis is ready."

# Symlink uploads directory
if [ -d /opt/postiz/uploads-default ] && [ ! -d /data/uploads/.initialized ]; then
    cp -rn /opt/postiz/uploads-default/* /data/uploads/ 2>/dev/null || true
    touch /data/uploads/.initialized
fi

# Change to Postiz directory
cd /opt/postiz || exit 1

# Run database migrations if needed
bashio::log.info "Running database migrations..."
npx prisma migrate deploy 2>&1 || bashio::log.warning "Migration failed or not needed"

# Start Postiz (the app listens on port 5000 internally)
# We use a simple socat proxy to map 4007 -> 5000 if needed,
# but first check if Postiz can be configured to listen on 4007
export PORT=4007

bashio::log.info "Launching Postiz on port ${PORT}..."
exec node /opt/postiz/dist/apps/backend/main.js &
BACKEND_PID=$!

# Start frontend (Next.js)
exec node /opt/postiz/dist/apps/frontend/server.js &
FRONTEND_PID=$!

# Wait for either to exit
wait -n ${BACKEND_PID} ${FRONTEND_PID}
