# postiz/rootfs/etc/s6-overlay/s6-rc.d/postgresql/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# PostgreSQL finish script
# ==============================================================================
bashio::log.info "PostgreSQL stopped."
