# postiz/rootfs/etc/s6-overlay/s6-rc.d/postgresql/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# PostgreSQL service
# ==============================================================================
bashio::log.info "Starting PostgreSQL..."

chown -R postgres:postgres /data/postgres /run/postgresql

exec su -s /bin/bash postgres -c "/usr/lib/postgresql/17/bin/postgres -D /data/postgres"
