# postiz/rootfs/etc/s6-overlay/s6-rc.d/postgresql/run
#!/command/with-contenv bashio
exec s6-setuidgid postgres /usr/lib/postgresql/17/bin/postgres -D /data/postgres
